# mopidy.service
[Unit]
Description=Mopidy        
After=network.target

[Service]
ExecStart=/usr/bin/mopidy

[Install]

