# config/launchd/com.bartek.ha-mcp-tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN"
   "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.bartek.ha-mcp-tunnel</string>

    <key>ProgramArguments</key>
    <array>
      <string>/usr/bin/ssh</string>
      <string>-N</string>
      <string>-o</string>
      <string>ServerAliveInterval=30</string>
      <string>-o</string>
      <string>ServerAliveCountMax=3</string>
      <string>-o</string>
      <string>ExitOnForwardFailure=yes</string>
      <string>-o</string>
      <string>StrictHostKeyChecking=accept-new</string>
      <string>-L</string>
      <string>8123:localhost:8123</string>
      <string>hass</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/bartek/Library/Logs/ha-mcp-tunnel.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/bartek/Library/Logs/ha-mcp-tunnel.err.log</string>
  </dict>
</plist>
