# systemd/gcalcli.service
[Unit]
Description=gcalcli remind

[Service]
Type=simple
ExecStart=/home/bata/dotfiles/gcal_remind.sh
