# systemd/gcalcli.timer
[Unit]
Description=gcalcli notification

[Timer]
OnBootSec=10min
OnUnitInactiveSec=30m
Unit=gcalcli.service

[Install]
WantedBy=multi-user.target
