# systemd/gcalcli.service
[Unit]
Description=gcalcli remind

[Service]
Type=simple
ExecStart=$HOME/dotfiles/systemd/gcal_remind.sh
