# systemd/mbsync.service
[Unit]
Description=check mails
RefuseManualStart=no
RefuseManualStop=yes

[Service]
Type=oneshot
ExecStart=/home/bata/dotfiles/systemd/get_mail.sh
